# export-lms-backups/export-lms-backups.crontab
0 16 * * * /bin/bash /opt/LMSbackup/lms.backup.sh
